# scripts/init/com.gitorchard.orchard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd unit for the orchard daemon (macOS).

  Install:
    cp scripts/init/com.gitorchard.orchard.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.gitorchard.orchard.plist

  Logs land in ~/.local/state/orchard/orchard.{out,err}.log so they
  coexist with the pidfile and any other daemon-owned state.

  RunAtLoad + KeepAlive give us "always running" behaviour. The daemon
  itself handles SIGTERM via signal.NotifyContext for clean shutdown
  during launchctl unload.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.gitorchard.orchard</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/orchard</string>
      <string>daemon</string>
      <string>start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/orchard.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/orchard.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
  </dict>
</plist>
